# tools/build/gaia/etc/systemd/system/gaia.service
[Unit]
Description=Cusp
Requires=network-online.target
After=network-online.target

[Service]
Restart=on-failure
User=cusp
Group=cusp
PermissionsStartOnly=true
ExecStart=/usr/bin/cusp node start --home=/etc/cusp
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target

